# config/opencode.service.example
[Unit]
Description=OpenCode API Server
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/<path-to-openbridge-workspace>
EnvironmentFile=%h/.config/openbridge/bridge.env
ExecStart=opencode serve --hostname 127.0.0.1 --port 4096
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
