# config/openbridge.service.example
[Unit]
Description=OpenBridge Telegram OpenCode Bridge
Wants=network-online.target
Wants=opencode.service
After=network-online.target opencode.service

[Service]
Type=simple
WorkingDirectory=%h/<path-to-openbridge-workspace>
EnvironmentFile=%h/.config/openbridge/bridge.env
ExecStart=%h/<path-to-openbridge-workspace>/.venv/bin/openbridge start --foreground
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
